# attitude_adjustment/files/etc/uci-defaults/luci-upnp
#!/bin/sh

/etc/init.d/miniupnpd enabled && {
	/etc/init.d/miniupnpd stop
	/etc/init.d/miniupnpd enable
}

rm -f /tmp/luci-indexcache
exit 0
